# root/etc/udev/rules.d/74-pr-disk.rules
KERNEL=="sd*", KERNELS=="1:0:0:0", SYMLINK+="removable%n", RUN+="/etc/init.d/hddtemp restart"

